# debian/huayra-apps-settings.postinst
#! /bin/sh
set -e

if [ "$1" = "configure" ]; then

  echo "Estableciendo preferencias de chromium"
  #TODO: hay que hacer refactor a dpkg-divert
  huayra-chromium-master-defaults
  update-rc.d fix-gdm3-autologin enable
fi

#DEBHELPER#
